# rootfs/etc/services.d/cache-manager/run
#!/usr/bin/with-contenv bash

if [[ -z "$CACHE_WATCH_DIR" ]]; then
    echo "CACHE_WATCH_DIR is not set. Disable cache."
    sleep 300
    exit
fi

MAX_CACHE_FILES=${MAX_CACHE_FILES:-2}

# Array to keep track of the LRU queue
declare -a CACHED_FILES

echo "Starting File RAM Cacher on $CACHE_WATCH_DIR..."

# Use process substitution < <(...) to avoid subshell variable issues
while read -r FILE
do
    # Only react to video files
    if [[ "$FILE" =~ \.(mkv|mp4|avi)$ ]]; then
        
        # Check if already in queue
        ALREADY_CACHED=0
        for cached in "${CACHED_FILES[@]}"; do
            if [[ "$cached" == "$FILE" ]]; then
                ALREADY_CACHED=1
                break
            fi
        done

        if [[ $ALREADY_CACHED -eq 0 ]]; then
            # IMPORTANT: Add to array BEFORE running vmtouch to prevent recursion loop
            CACHED_FILES+=("$FILE")
            echo "New File detected: $FILE. Caching..."
            
            # Cache in background
            vmtouch -t "$FILE" > /dev/null 2>&1 &

            # LRU Eviction Logic
            if [[ ${#CACHED_FILES[@]} -gt $MAX_CACHE_FILES ]]; then
                OLDEST_FILE="${CACHED_FILES[0]}"
                echo "Cache full. Evicting: $OLDEST_FILE"
                vmtouch -e "$OLDEST_FILE" > /dev/null 2>&1
                CACHED_FILES=("${CACHED_FILES[@]:1}")
            fi
        fi
    fi
# Note the syntax change here to fix the subshell issue:
done < <(inotifywait -m -e open --format "%w%f" -r "$CACHE_WATCH_DIR")
