# rootfs/etc/services.d/cache-manager/run
#!/usr/bin/with-contenv bash

# --- Configuration ---
if [[ -z "$CACHE_WATCH_DIR" ]]; then
    echo "[ERROR] CACHE_WATCH_DIR is not set. Script exiting."
    sleep 60
    exit 1
fi

# Associative array to store the last time a file was touched
declare -A LAST_TOUCH_TIME
COOLDOWN=30 # Seconds to ignore a file after the first 'open'

echo "--- Simple Media RAM Cacher Started ---"
echo "Watching: $CACHE_WATCH_DIR"
echo "Strategy: Cache on Open / Automatic Linux Eviction"

# Monitor for 'open' events recursively
while read -r FILE
do
    # 1. Only react to video files
    if [[ "$FILE" =~ \.(mkv|mp4|avi|ts|m4v|mks)$ ]]; then
        
        NOW=$(date +%s)
        LAST=${LAST_TOUCH_TIME["$FILE"]:-0}

        # 2. Loop Prevention (Cooldown)
        # If we touched this file less than 30 seconds ago, skip it.
        # This prevents the script from reacting to vmtouch's own 'open' action.
        if (( NOW - LAST < COOLDOWN )); then
            continue
        fi

        # 3. Update the timestamp and Cache the file
        LAST_TOUCH_TIME["$FILE"]=$NOW
        
        echo "[CACHE] Loading: $FILE"
        
        # -t touches the file into RAM.
        # Running in background (&) so the script keeps listening for new files.
        vmtouch -t "$FILE" > /dev/null 2>&1 &

        # 4. Prevent memory leak in the script itself
        if [[ ${#LAST_TOUCH_TIME[@]} -gt 10 ]]; then
            unset LAST_TOUCH_TIME
            declare -A LAST_TOUCH_TIME
        fi
    fi

# Use process substitution to feed inotifywait into the loop
done < <(inotifywait -m -e open --format "%w%f" -r "$CACHE_WATCH_DIR" 2>/dev/null)
