# rootfs/etc/services.d/cache-manager/run
#!/usr/bin/with-contenv bash
exec s6-setuidgid root env /usr/local/bin/cache-manager
